# Linux Shell Bash/dump_database.timer
[Unit]
Description=Run my service every 10mins

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
