# NVIDIA/TD-9-11-22/dump_database.service
[Unit]
Description=Start database backup by executing dump_database

[Service]
WorkingDirectory=/home
Type=simple
ExecStart=/home/dump_database.sh
